# deploy/homelab-cluster-agent.service
[Unit]
Description=Homelab cluster agent (startup/shutdown webhook)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# kubectl (invoked by startup.sh's phase scripts) resolves its default
# kubeconfig from $HOME/.kube/config. systemd doesn't set $HOME for this
# unit on its own, so without this, kubectl silently fails to find any
# config and wake-control-plane.sh's "wait for API server" loop retries
# forever with no visible error (its kubectl call is `&>/dev/null`).
Environment=HOME=/root
EnvironmentFile=/etc/homelab-cluster-agent.env
ExecStart=/usr/local/bin/homelab-cluster-agent
Restart=on-failure
RestartSec=5

# Hardening — no network overhead beyond the process itself, no container runtime.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
